# deploy/systemd/zhifei-bid.service
[Unit]
Description=Qingtian Bid System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/zhifei/ZhiFei_BizSystem
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/zhifei-bid.env
ExecStart=/opt/zhifei/ZhiFei_BizSystem/.venv/bin/python -m app.main --no-browser
Restart=always
RestartSec=3
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
